# systemd/speak-api.service
[Unit]
Description=Speak Small API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/speak-small
EnvironmentFile=/opt/speak-small/.env
Environment=ADDR=127.0.0.1:18080
ExecStart=/opt/speak-small/bin/theme-server
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
